# citadel/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/citadel
drwxr-xr-x	root/root	var/
drwxr-xr-x	citadel/citadel	var/citadel/
-rw-r--r--	citadel/citadel	var/citadel/README.txt
-rwxr-xr-x	citadel/citadel	var/citadel/aidepost
-rwxr-xr-x	citadel/citadel	var/citadel/base64
drwxr-xr-x	citadel/citadel	var/citadel/bio/
-rwxr-xr-x	citadel/citadel	var/citadel/chkpw
-rwxr-xr-x	citadel/citadel	var/citadel/chkpwd
-rw-r--r--	citadel/citadel	var/citadel/citadel_urlshorteners.rc
-rwxr-xr-x	citadel/citadel	var/citadel/citmail
-rwxr-xr-x	citadel/citadel	var/citadel/citserver
-rwxr-xr-x	citadel/citadel	var/citadel/ctdlmigrate
-rwxr-xr-x	citadel/citadel	var/citadel/database_cleanup.sh
drwxr-xr-x	citadel/citadel	var/citadel/docs/
-rw-r--r--	citadel/citadel	var/citadel/docs/PAM.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/README-FIRST.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/binaries.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/build.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/citadel_thread_io.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/citadelapi.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/databaselayout.md
-rw-r--r--	citadel/citadel	var/citadel/docs/delivery-list.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/developers.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/logging.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/netconfigs.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/package-setup.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/smtpclient.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/test_suite.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/views.txt
-rw-r--r--	citadel/citadel	var/citadel/docs/welcomemail.html
-rw-r--r--	citadel/citadel	var/citadel/docs/welcomemail.txt
drwxr-xr-x	citadel/citadel	var/citadel/files/
drwxr-xr-x	citadel/citadel	var/citadel/images/
drwxr-xr-x	citadel/citadel	var/citadel/info/
drwxr-xr-x	citadel/citadel	var/citadel/locale/
drwxr-xr-x	citadel/citadel	var/citadel/locale/ar/
drwxr-xr-x	citadel/citadel	var/citadel/locale/ar/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/ar/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/cs/
drwxr-xr-x	citadel/citadel	var/citadel/locale/cs/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/cs/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/de/
drwxr-xr-x	citadel/citadel	var/citadel/locale/de/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/de/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/en_GB/
drwxr-xr-x	citadel/citadel	var/citadel/locale/en_GB/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/en_GB/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/es/
drwxr-xr-x	citadel/citadel	var/citadel/locale/es/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/es/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/fi/
drwxr-xr-x	citadel/citadel	var/citadel/locale/fi/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/fi/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/fr/
drwxr-xr-x	citadel/citadel	var/citadel/locale/fr/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/fr/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/hu/
drwxr-xr-x	citadel/citadel	var/citadel/locale/hu/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/hu/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/it/
drwxr-xr-x	citadel/citadel	var/citadel/locale/it/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/it/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/ja/
drwxr-xr-x	citadel/citadel	var/citadel/locale/ja/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/ja/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/kk/
drwxr-xr-x	citadel/citadel	var/citadel/locale/kk/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/kk/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/ko/
drwxr-xr-x	citadel/citadel	var/citadel/locale/ko/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/ko/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/pl/
drwxr-xr-x	citadel/citadel	var/citadel/locale/pl/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/pl/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/pt/
drwxr-xr-x	citadel/citadel	var/citadel/locale/pt/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/pt/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/pt_BR/
drwxr-xr-x	citadel/citadel	var/citadel/locale/pt_BR/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/pt_BR/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/ru/
drwxr-xr-x	citadel/citadel	var/citadel/locale/ru/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/ru/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/sv/
drwxr-xr-x	citadel/citadel	var/citadel/locale/sv/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/sv/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/locale/vi/
drwxr-xr-x	citadel/citadel	var/citadel/locale/vi/LC_MESSAGES/
-rwxr-xr-x	citadel/citadel	var/citadel/locale/vi/LC_MESSAGES/citadel-setup.mo
drwxr-xr-x	citadel/citadel	var/citadel/messages/
-rw-r--r--	citadel/citadel	var/citadel/messages/changepw
-rw-r--r--	citadel/citadel	var/citadel/messages/goodbye
-rw-r--r--	citadel/citadel	var/citadel/messages/hello
-rw-r--r--	citadel/citadel	var/citadel/messages/newuser
-rw-r--r--	citadel/citadel	var/citadel/messages/register
-rw-r--r--	citadel/citadel	var/citadel/messages/roomaccess
-rw-r--r--	citadel/citadel	var/citadel/messages/unlisted
-rwxr-xr-x	citadel/citadel	var/citadel/migrate_aliases.sh
-rwxr-xr-x	citadel/citadel	var/citadel/msgform
drwxr-xr-x	citadel/citadel	var/citadel/network/
-rw-r--r--	citadel/citadel	var/citadel/network/mail.aliases
drwxr-xr-x	citadel/citadel	var/citadel/network/spoolin/
drwxr-xr-x	citadel/citadel	var/citadel/network/spoolout/
-rw-r--r--	citadel/citadel	var/citadel/public_clients
-rwxr-xr-x	citadel/citadel	var/citadel/sendcommand
-rwxr-xr-x	citadel/citadel	var/citadel/setup
drwxr-xr-x	citadel/citadel	var/citadel/techdoc/
drwxr-xr-x	citadel/citadel	var/citadel/unstripped/
-rwxr-xr-x	citadel/citadel	var/citadel/unstripped/citserver
drwxr-xr-x	citadel/citadel	var/citadel/userpics/
