# fprintd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/fprint
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/security/
-rw-r--r--	root/root	lib/security/pam_fprintd.a
-rwxr-xr-x	root/root	lib/security/pam_fprintd.la
-rwxr-xr-x	root/root	lib/security/pam_fprintd.so
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/fprintd-delete
-rwxr-xr-x	root/root	usr/bin/fprintd-enroll
-rwxr-xr-x	root/root	usr/bin/fprintd-list
-rwxr-xr-x	root/root	usr/bin/fprintd-verify
drwxr-xr-x	root/root	usr/etc/
-rw-r--r--	root/root	usr/etc/fprintd.conf
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/fprintd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/interfaces/
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/net.reactivated.Fprint.Device.xml
-rw-r--r--	root/root	usr/share/dbus-1/interfaces/net.reactivated.Fprint.Manager.xml
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/net.reactivated.Fprint.service
drwxr-xr-x	root/root	usr/share/dbus-1/system.d/
-rw-r--r--	root/root	usr/share/dbus-1/system.d/net.reactivated.Fprint.conf
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/fprintd.1.gz
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/net.reactivated.fprint.device.policy
drwxr-xr-x	root/root	usr/var/
drwxr-xr-x	root/root	usr/var/lib/
drwxr-xr-x	root/root	usr/var/lib/fprint/
