# opensmtpd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/opensmtpd
-rw-r--r--	root/root	etc/smtpd.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/smtp
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/opensmtpd/
-rwxr-xr-x	root/root	usr/lib/opensmtpd/encrypt
-r-sr-xr-x	root/root	usr/lib/opensmtpd/lockspool
-rwxr-xr-x	root/root	usr/lib/opensmtpd/mail.lmtp
-rwxr-xr-x	root/root	usr/lib/opensmtpd/mail.local
-rwxr-xr-x	root/root	usr/lib/opensmtpd/mail.maildir
-rwxr-xr-x	root/root	usr/lib/opensmtpd/mail.mboxfile
-rwxr-xr-x	root/root	usr/lib/opensmtpd/mail.mda
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/lockspool.1.gz
-rw-r--r--	root/root	usr/man/man1/smtp.1.gz
drwxr-xr-x	root/root	usr/man/man5/
-rw-r--r--	root/root	usr/man/man5/aliases.5.gz
-rw-r--r--	root/root	usr/man/man5/forward.5.gz
-rw-r--r--	root/root	usr/man/man5/smtpd.conf.5.gz
-rw-r--r--	root/root	usr/man/man5/table.5.gz
drwxr-xr-x	root/root	usr/man/man7/
-rw-r--r--	root/root	usr/man/man7/smtpd-filters.7.gz
drwxr-xr-x	root/root	usr/man/man8/
-rw-r--r--	root/root	usr/man/man8/mail.lmtp.8.gz
-rw-r--r--	root/root	usr/man/man8/mail.local.8.gz
-rw-r--r--	root/root	usr/man/man8/mail.maildir.8.gz
-rw-r--r--	root/root	usr/man/man8/mail.mboxfile.8.gz
-rw-r--r--	root/root	usr/man/man8/mail.mda.8.gz
-rw-r--r--	root/root	usr/man/man8/makemap.8.gz
-rw-r--r--	root/root	usr/man/man8/newaliases.8.gz
-rw-r--r--	root/root	usr/man/man8/sendmail.8.gz
-rw-r--r--	root/root	usr/man/man8/smtpctl.8.gz
-rw-r--r--	root/root	usr/man/man8/smtpd.8.gz
drwxr-xr-x	root/root	usr/sbin/
lrwxrwxrwx	root/root	usr/sbin/sendmail -> /usr/sbin/smtpctl
-rwxr-sr-x	root/mail	usr/sbin/smtpctl
-rwxr-xr-x	root/root	usr/sbin/smtpd
