# clamsmtp/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/clamav/
-rw-r--r--	root/root	etc/clamav/clamsmtpd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/clamsmtp
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/clamsmtpd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/clamsmtpd.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/clamsmtpd.8.gz
