# spamassassin/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/spamd
drwxr-xr-x	root/root	etc/spamassassin/
-rw-r--r--	root/root	etc/spamassassin/init.pre
-rw-r--r--	root/root	etc/spamassassin/local.cf
-rw-r--r--	root/root	etc/spamassassin/v310.pre
-rw-r--r--	root/root	etc/spamassassin/v312.pre
-rw-r--r--	root/root	etc/spamassassin/v320.pre
-rw-r--r--	root/root	etc/spamassassin/v330.pre
-rw-r--r--	root/root	etc/spamassassin/v340.pre
-rw-r--r--	root/root	etc/spamassassin/v341.pre
-rw-r--r--	root/root	etc/spamassassin/v342.pre
-rw-r--r--	root/root	etc/spamassassin/v343.pre
-rw-r--r--	root/root	etc/spamassassin/v400.pre
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-r-xr-xr-x	root/root	usr/bin/sa-awl
-r-xr-xr-x	root/root	usr/bin/sa-check_spamd
-r-xr-xr-x	root/root	usr/bin/sa-compile
-r-xr-xr-x	root/root	usr/bin/sa-learn
-r-xr-xr-x	root/root	usr/bin/sa-update
-r-xr-xr-x	root/root	usr/bin/spamassassin
-r-xr-xr-x	root/root	usr/bin/spamc
-r-xr-xr-x	root/root	usr/bin/spamd
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/perl5/
drwxr-xr-x	root/root	usr/lib/perl5/5.36/
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/AICache.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/ArchiveIterator.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/AsyncLoop.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/AutoWelcomelist.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Bayes.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Bayes/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Bayes/CombineChi.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Bayes/CombineNaiveBayes.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/BDB.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/DBM.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/MySQL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/PgSQL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/Redis.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/SDBM.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/BayesStore/SQL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Client.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Conf.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Conf/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Conf/LDAP.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Conf/Parser.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Conf/SQL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Constants.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/DBBasedAddrList.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Dns.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/DnsResolver.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/GeoDB.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/HTML.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Locales.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Locker.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Locker/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Locker/Flock.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Locker/UnixNFSSafe.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Locker/Win32.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Logger.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Logger/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Logger/File.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Logger/Stderr.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Logger/Syslog.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/MailingList.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Message.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Message/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Message/Metadata.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Message/Metadata/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Message/Metadata/Received.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Message/Node.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/NetSet.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/PerMsgLearner.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/PerMsgStatus.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/PersistentAddrList.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/ASN.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/AWL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/AccessDB.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/AntiVirus.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/AskDNS.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/AuthRes.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/AutoLearnThreshold.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Bayes.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/BodyEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/BodyRuleBaseExtractor.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Check.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/DCC.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/DKIM.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/DMARC.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/DNSEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/DecodeShortURLs.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/ExtractText.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/FreeMail.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/FromNameSpoof.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/HTMLEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/HTTPSMismatch.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/HashBL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/HeaderEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/ImageInfo.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/MIMEEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/MIMEHeader.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/OLEVBMacro.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/OneLineBodyRuleType.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/PDFInfo.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/PhishTag.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Phishing.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Pyzor.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Razor2.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/RelayCountry.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/RelayEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/ReplaceTags.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/ResourceLimits.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Reuse.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Rule2XSBody.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/SPF.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Shortcircuit.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/SpamCop.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/Test.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/TextCat.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/TxRep.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/URIDNSBL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/URIDetail.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/URIEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/URILocalBL.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/VBounce.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/WLBLEval.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Plugin/WelcomeListSubject.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/PluginHandler.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/RegistryBoundaries.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Reporter.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/SQLBasedAddrList.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/SpamdForkScaling.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/SubProcBackChannel.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Timeout.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Util.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Util/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Util/DependencyInfo.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Util/Progress.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Util/ScopedTimer.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Util/TieOneStringHash.pm
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/Mail/SpamAssassin/Util/TinyRedis.pm
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/linux-thread-multi/
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/linux-thread-multi/auto/
drwxr-xr-x	root/root	usr/lib/perl5/site_perl/5.36/linux-thread-multi/auto/Mail/
-r--r--r--	root/root	usr/lib/perl5/site_perl/5.36/spamassassin-run.pod
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-r--r--r--	root/root	usr/share/man/man1/sa-awl.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/sa-check_spamd.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/sa-compile.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/sa-learn.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/sa-update.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/spamassassin-run.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/spamassassin.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/spamc.1pm.gz
-r--r--r--	root/root	usr/share/man/man1/spamd.1pm.gz
drwxr-xr-x	root/root	usr/share/man/man3/
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::AICache.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::ArchiveIterator.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::AsyncLoop.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::AutoWelcomelist.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Bayes.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::BayesStore.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::BayesStore::BDB.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::BayesStore::MySQL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::BayesStore::PgSQL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::BayesStore::Redis.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::BayesStore::SQL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Client.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Conf.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Conf::LDAP.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Conf::Parser.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Conf::SQL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::DnsResolver.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::GeoDB.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Logger.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Logger::File.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Logger::Stderr.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Logger::Syslog.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Message.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Message::Metadata.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Message::Node.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::PerMsgLearner.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::PerMsgStatus.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::PersistentAddrList.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::ASN.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::AWL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::AccessDB.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::AntiVirus.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::AskDNS.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::AuthRes.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::AutoLearnThreshold.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Bayes.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::BodyRuleBaseExtractor.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Check.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::DCC.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::DKIM.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::DMARC.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::DNSEval.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::DecodeShortURLs.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::ExtractText.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::FreeMail.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::FromNameSpoof.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::HashBL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::MIMEEval.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::MIMEHeader.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::OLEVBMacro.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::OneLineBodyRuleType.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::PDFInfo.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::PhishTag.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Phishing.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Pyzor.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Razor2.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::RelayCountry.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::ReplaceTags.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::ResourceLimits.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Reuse.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Rule2XSBody.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::SPF.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Shortcircuit.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::SpamCop.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::Test.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::TextCat.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::TxRep.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::URIDNSBL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::URIDetail.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::URILocalBL.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::VBounce.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Plugin::WelcomeListSubject.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::PluginHandler.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::RegistryBoundaries.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::SQLBasedAddrList.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::SubProcBackChannel.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Timeout.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Util.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Util::DependencyInfo.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/Mail::SpamAssassin::Util::Progress.3pm.gz
-r--r--r--	root/root	usr/share/man/man3/spamassassin-run.3pm.gz
drwxr-xr-x	root/root	usr/share/spamassassin/
-rw-r--r--	root/root	usr/share/spamassassin/languages
-rw-r--r--	root/root	usr/share/spamassassin/sa-update-pubkey.txt
-rw-r--r--	root/root	usr/share/spamassassin/user_prefs.template
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	spamd/spamd	var/lib/spamassassin/
